# supervisord-exim4.conf
[program:exim4]
command=/start-exim4.sh
numprocs=1
autostart=true
autorestart=true
